# com.localtime.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.localtime.desktop</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>FIREBASE_PROJECT_ID</key>
        <string>localtime-cloud-wm</string>
        <key>FIREBASE_API_KEY</key>
        <string>AIzaSyAsrp1py46ECzlF_eE6LJG7-FzdBBPmxZg</string>
    </dict>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/wilbmoffitt/Desktop/LocalTime/desktop_agent.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/wilbmoffitt/Desktop/LocalTime/agent_out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/wilbmoffitt/Desktop/LocalTime/agent_err.log</string>
</dict>
</plist>
